# ioc-checker.service
[Unit]
Description=IOC Checker HTTP Service
After=network.target

[Service]
WorkingDirectory=/opt/ioc-checker-playwright
ExecStart=/usr/bin/env python -m hypercorn ioc_checker.main:app --bind 0.0.0.0:8000
Restart=on-failure

[Install]
WantedBy=multi-user.target
